# rootfs/etc/s6-overlay/s6-rc.d/svc-cron/finish
#!/usr/bin/with-contenv bash

echo "Caught SIGTERM signal!"

kill -TERM "$(cat /var/run/crond.pid)" 2>/dev/null
