# rootfs/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/usr/bin/with-contenv bash

echo "Starting cron daemon..."
exec s6-setuidgid root /usr/sbin/crond -f -l 2
