# rootfs/etc/s6-overlay/s6-rc.d/svc-housekeeper/finish
#!/usr/bin/with-contenv bash

echo "Caught SIGTERM signal!"

kill -TERM "$(cat /var/run/housekeeperd.pid)" 2>/dev/null
kill -TERM "$(cat /var/run/backupd.pid)" 2>/dev/null